# YAML/101-ansible/galera-ha-cluster/roles/galera/templates/galera-agent-rw.socket.j2
[Socket]
ListenStream={{ galera_agent_port_rw }}
Accept=false

[Install]
WantedBy=sockets.target
